# crontab
0 0 * * * /opt/minecraft/backup/mc-backup.rb 2>> /opt/minecraft/backup/log/backup.log
